# deploy/practice-field-scheduler.service
[Unit]
Description=Field Scheduler Next.js Application
After=network.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/opt/practice-field-scheduler
Environment=NODE_ENV=production
Environment=PORT=3001
EnvironmentFile=/opt/practice-field-scheduler/.env

ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 